# content/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=/usr/sbin/apachectl -D FOREGROUND

[program:mariadb]
command=/usr/bin/mysqld_safe

[program:icinga2]
command=icinga2 daemon
